# root/etc/s6-overlay/s6-rc.d/init-calibre-web-config/run
#!/usr/bin/with-contenv bash

# create symlinks for imagemagick policy.xml
rm -rf /etc/ImageMagick-6/policy.xml
ln -s /defaults/policy.xml /etc/ImageMagick-6/policy.xml

# create symlinks for database and log
if [[ ! -f /config/app.db ]] && [[ -f /app/calibre-web/app.db ]]; then
    cp /defaults/app.db /config/app.db
fi
rm -rf /app/calibre-web/app.db
ln -s /config/app.db /app/calibre-web/app.db

# create symlinks for log
if [[ ! -f /config/calibre-web.log ]]; then
    touch /config/calibre-web.log
fi
rm -rf /app/calibre-web/calibre-web.log
ln -s /config/calibre-web.log /app/calibre-web/calibre-web.log

# create Google drive client_secrets.json file
if [[ ! -f /config/client_secrets.json ]]; then
    echo "{}" > /config/client_secrets.json
fi
rm -rf /app/calibre-web/client_secrets.json
ln -s /config/client_secrets.json /app/calibre-web/client_secrets.json

# create Google drive symlinks for database
if [[ ! -f /config/gdrive.db ]] && [[ -f /app/calibre-web/gdrive.db ]]; then
    cp /app/calibre-web/gdrive.db /config/gdrive.db
fi
rm -rf /app/calibre-web/gdrive.db
ln -s /config/gdrive.db /app/calibre-web/gdrive.db

# check if kepubify is present and if so make executable
if [[ -f /usr/bin/kepubify ]] && [[ ! -x /usr/bin/kepubify ]]; then
    chmod +x /usr/bin/kepubify
fi

# Symlink logfiles
if [[ ! -L "/app/calibre-web/calibre-web.log" ]]; then
    ln -s /config/calibre-web.log /app/calibre-web/calibre-web.log
fi

# permissions
lsiown -R abc:abc \
    /config
