# root/etc/services.d/calibre-web/run
#!/usr/bin/with-contenv bash

cd /app/calibre-web || exit

exec \
	s6-setuidgid abc /usr/bin/gunicorn \
	-c /app/calibre-web/gunicorn.conf.py wsgi:web.app
