# root/etc/services.d/calibre-web/run
#!/usr/bin/with-contenv bash

cd /app/calibre-web || exit

exec \
    s6-setuidgid abc python3 /app/calibre-web/cps.py
